# supervisord/bandtools.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:nginx]
command=nginx -c /app/nginx.conf -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:nextjs]
command=node /app/setlist-helper/.next/standalone/server.js
directory=/app/setlist-helper/.next/standalone
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PORT="3000",HOSTNAME="0.0.0.0"
priority=20

[program:setloader]
command=/opt/setloader-venv/bin/uvicorn server_simple_db:app --host 127.0.0.1 --port 8002
directory=/app/setloader
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30

[program:flyers]
command=/opt/flyers-venv/bin/uvicorn bridge.server:app --host 127.0.0.1 --port 8080
directory=/app/gig-flyers
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=ROOT_PATH="/flyers",GIG_OUTPUT_DIR="/data/flyers-output",GIG_CALENDAR_CACHE_PATH="/data/cache/calendar.json",BRIDGE_POLL_SECONDS="86400"
priority=40
